# services/cookie-manager/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1280x720x24
autorestart=true
priority=10

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":99"
autorestart=true
priority=20

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -nopw -shared -rfbport 5900
autorestart=true
priority=30

[program:novnc]
command=/usr/bin/websockify --web /usr/share/novnc 6901 localhost:5900
autorestart=true
priority=40

[program:chromium-gui]
command=/usr/local/bin/chromium --display=:99 --no-sandbox --disable-gpu --user-data-dir=/data/chrome-profile --start-maximized https://www.youtube.com
autorestart=false
startsecs=5
priority=50

[program:health-api]
command=/usr/local/bin/uvicorn main:app --host 0.0.0.0 --port 6902
directory=/app
environment=DISPLAY=":99"
autorestart=true
priority=60
